# deploy/dashboard.service
[Unit]
Description=Personal Dashboard (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
# Ajusta User se não correres como root (ex.: um utilizador 'dashboard').
User=root
WorkingDirectory=/opt/dashboard
# uvicorn só no localhost — exposto ao mundo apenas via Nginx (porta 80).
ExecStart=/opt/dashboard/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
